# tmp_fs/etc/init.d/rc.local
#!/bin/sh
#add product specified script here
